# docker-novnc/conf.d/xterm.conf
[program:steam]
command=steam
autorestart=true
